# .agents/skills/wecom-cli-bot/assets/wecom-cli-bots-template/supervisor/launchd/com.local.wecom-cli-bot.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.wecom-cli-bot.BOT_NAME</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/ABSOLUTE/PATH/wecom-cli-bots/dist/main.js</string>
    <string>--bot</string>
    <string>BOT_NAME</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/wecom-cli-bots</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
