# .agents/skills/wecom-cli-bot/assets/wecom-cli-bots-template/supervisor/systemd/wecom-cli-bot@.service
[Unit]
Description=WeCom CLI Bot %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/wecom-cli-bots
ExecStart=/usr/bin/node /opt/wecom-cli-bots/dist/main.js --bot %i
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
